# supervisord.conf
[program:webserver]
command=gunicorn -b '127.0.0.1:{{ settings.SERVER_PORT }}' --workers 3 core.wsgi:application --error-logfile /home/cityparks/cityparksphila.org/logs/gunicorn/error.log
